# recipes-app/adi/adrv906x-files/files/bootsuccess.service
[Unit]
Description=boot success setting service

[Service]
ExecStart=/usr/bin/optee_app_boot_successful

[Install]
WantedBy=multi-user.target
